# src/macrandrd
#!/bin/ksh
#
# $OpenBSD: macrandr,v 1.2 2019/08/06 00:02:12 rpe Exp $

daemon="/usr/sbin/macrandr -D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1